# roles/socat_proxy/templates/socat-proxy.socket.j2
[Unit]
Description=Listening Socket for socat proxy (443)

[Socket]
ListenStream=443
Accept=yes

[Install]
WantedBy=sockets.target
