# roles/socat_proxy/templates/wol@.service.j2
[Unit]
Description=Sending Wake on LAN packet to %i
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wakeonlan -i 255.255.255.255 %i
